# .build/packages/lite-nas_0.1.0_arm64.root/usr/libexec/lite-nas/configs/systemd/system/lite-nas-system-metrics.service
[Unit]
Description=LiteNAS system metrics service
After=network-online.target nats-server.service
Wants=network-online.target
Requires=nats-server.service

[Service]
Type=simple
User=@SYSTEM_METRICS_RUNTIME_USER@
Group=@SYSTEM_METRICS_RUNTIME_GROUP@
SupplementaryGroups=@SYSTEM_METRICS_CONFIG_GROUP@
ExecStart=@SYSTEM_METRICS_BINARY@
Restart=on-failure
RestartSec=5s
WorkingDirectory=/
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectProc=invisible
ProtectSystem=strict
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=
UMask=0027
ReadOnlyPaths=@SYSTEM_METRICS_CONFIG_DIR@ /proc/stat /proc/meminfo
ReadWritePaths=@SYSTEM_METRICS_LOG_FILE@
StandardOutput=null
StandardError=journal

[Install]
WantedBy=multi-user.target
